# scripts/ai-code-review-worker.service
[Unit]
Description=AI Code Review System - Queue Worker
Documentation=https://github.com/your-org/ai-code-review-system
After=network.target redis.service mysql.service
Wants=redis.service mysql.service
Requires=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/ai-code-review-system
ExecStart=/usr/bin/tsx /opt/ai-code-review-system/scripts/worker.ts
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=35
Restart=always
RestartSec=10
StartLimitInterval=60
StartLimitBurst=3

# 环境变量
Environment=NODE_ENV=production
Environment=LOG_LEVEL=info
Environment=WORKER_MAX_CONCURRENCY=10
Environment=WORKER_POLL_INTERVAL=1000
Environment=WORKER_TASK_TIMEOUT=600000
Environment=WORKER_SHUTDOWN_TIMEOUT=30000
Environment=WORKER_HEALTH_CHECK_INTERVAL=30000
Environment=WORKER_ENABLE_METRICS=true

# 从环境文件加载配置
EnvironmentFile=-/opt/ai-code-review-system/.env
EnvironmentFile=-/etc/ai-code-review-system/worker.env

# 安全设置
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/ai-code-review-system/logs
ReadWritePaths=/tmp

# 资源限制
LimitNOFILE=65536
LimitNPROC=4096
MemoryMax=2G
CPUQuota=200%

# 日志设置
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ai-code-review-worker

[Install]
WantedBy=multi-user.target